# guest/overlay/etc/profile.d/cosim.sh
#!/bin/sh
# Ensure cosim tools are in PATH
export PATH="/usr/local/bin:$PATH"
